# theos-recovery.service
[Unit]
Description=theOS Recovery Environment

[Service]
ExecStart=/usr/bin/theos-recovery

[Install]
WantedBy=emergency.target
